# deploy/systemd/earthquake-alert.service
[Unit]
Description=Earthquake Alert Bark subscription service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/earthquake-alert
EnvironmentFile=-/opt/earthquake-alert/.env
ExecStart=/opt/earthquake-alert/earthquake-alert-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
